# roles/backup/templates/restic-backup.service.j2
[Unit]
Description=Run Restic Backup
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-backup.sh
